# deploy/supervisor/flockiq-gunicorn.conf
[program:flockiq-gunicorn]
command=/www/wwwroot/flockiq/venv/bin/gunicorn
    config.wsgi:application
    --bind 127.0.0.1:8002
    --workers 4
    --worker-class sync
    --timeout 120
    --keepalive 5
    --max-requests 1000
    --max-requests-jitter 50
    --log-level warning
    --access-logfile /www/wwwroot/flockiq/logs/gunicorn-access.log
    --error-logfile /www/wwwroot/flockiq/logs/gunicorn-error.log
directory=/www/wwwroot/flockiq
user=www
autostart=true
autorestart=true
startsecs=5
stopasgroup=true
killasgroup=true
environment=DJANGO_SETTINGS_MODULE="config.settings.production"
stderr_logfile=/www/wwwroot/flockiq/logs/gunicorn.err.log
stdout_logfile=/www/wwwroot/flockiq/logs/gunicorn.out.log
